# docker/supervisord.conf
# =============================================================================
# supervisord.conf — Process manager configuration
# =============================================================================
# Supervisord keeps Nginx and PHP-FPM running inside a single container.
# If either process dies, Supervisord restarts it automatically.
# This is the standard pattern for single-container PHP+Nginx deployments.

[supervisord]
nodaemon=true                   ; Run supervisord in the foreground (required for Docker)
user=root                       ; Must be root to manage both www-data processes
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

# ─────────────────────────────────────────────────────────────────────────────
# Nginx web server
# ─────────────────────────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ─────────────────────────────────────────────────────────────────────────────
# PHP-FPM application server
# ─────────────────────────────────────────────────────────────────────────────
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
priority=5                       ; Start before nginx so the socket is ready
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ─────────────────────────────────────────────────────────────────────────────
# Laravel Queue Worker (database driver)
# ─────────────────────────────────────────────────────────────────────────────
# Your app uses QUEUE_CONNECTION=database. This worker processes queued jobs
# (e.g. notifications). Remove this block if you don't use queues.
[program:queue-worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
user=www-data
autostart=true
autorestart=true
priority=20
numprocs=1
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopwaitsecs=3600
